# Middle_Layer/Middle_Layer_v3-0-1/EMU/DCS/HighVoltage/component/source/CMS_CSCfw_HV_CC/special_stuff/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
##/home/dcs/bin/dcscontrol_watch_dog &
/home/fast/bin/dns_wrap &
mkdir /mnt/central_installation
mkdir /mnt/central_repository
mount -t nfs pcmsucr3:/central_dcs_installation /mnt/central_installation &
mount -t nfs pcmsucr3:/central_dcs_repositoty_csc /mnt/central_repository &